# fs/etc/local.d/initialscript.start
#!/bin/bash

INITIAL_PATH=$(pwd)

cd /etc/init.d
for n in `seq 1 6`; do ln -s agetty agetty.tty$n; rc-config add agetty.tty$n default; done

cd $INITIAL_PATH
mv /etc/local.d/initialscript.start /etc/initialscript.start.used